# rootfs/etc/services.d/dnsmasq/run
#!/usr/bin/with-contenv sh
exec dnsmasq --keep-in-foreground --conf-file=/etc/dnsmasq.conf